# files0/files/etc/crontabs/root
# 每周3、6凌晨5点重启路由器系统，延时70秒
0 5  3,6 sleep 70 && touch /etc/banner && reboot
